# elk/stuff/supervisord.conf
[supervisord]
nodaemon=true

[program:httpd]
command=/usr/sbin/httpd -DFOREGROUND

[program:filebeat]
command=/usr/bin/filebeat -e